# ninjapack/root/etc/profile.d/disable-pairing.sh
export NINJA_SPHERE_CLIENT_ARGS=${NINJA_SPHERE_CLIENT_ARGS}${NINJA_SPHERE_CLIENT_ARGS:+ }--noCloud
